# data/storaged.service.in
[Unit]
Description=Storaged
Documentation=man:storaged(8)

[Service]
Type=dbus
BusName=org.storaged.Storaged
ExecStart=@storagedprivdir@/storaged --no-debug
