# data/udisks2.service.in
[Unit]
Description=Storage Daemon

[Service]
Type=dbus
BusName=org.freedesktop.UDisks2
ExecStart=@udisksdprivdir@/udisksd --no-debug
