# data/udisks2.service.in
[Unit]
Description=Disk Manager
Documentation=man:udisks(8)

[Service]
Type=dbus
BusName=org.freedesktop.UDisks2
ExecStart=@udisksdprivdir@/udisksd
KillSignal=SIGINT

[Install]
WantedBy=graphical.target
