# src/99-devkit-disks.rules
# pass all events to the DeviceKit disks daemon
SUBSYSTEM=="block", RUN+="socket:/org/freedesktop/devicekit/disks/udev_event"
